# smut@mdb.service
[Unit]
Description=Simple Mender Update Tool (MDB)
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
ExecStart=/usr/local/bin/smut \
    --redis-addr=192.168.7.1:6379 \
    --update-key=mender/update/mdb/url \
    --checksum-key=mender/update/mdb/checksum \
    --failure-key=mender/update/last-failure \
    --download-dir=/var/lib/mender/download \
    --update-type=non-blocking
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
